# quizly.service
[Unit]
Description=Quizly Frontend Service
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/quizly-frontend
ExecStart=/usr/bin/serve -s build -l 3000
Restart=always
Environment=NODE_ENV=production

[Install]
WantedBy=multi-user.target